# com.vcommits.ghidorah.node3.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vcommits.ghidorah.node3</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Adjust path to your iMac's shell script -->
        <string>/Users/vcommits/ghidorah_node_core/start_imac_node.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ghidorah_node3.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ghidorah_node3_err.log</string>
</dict>
</plist>
